# netbird-exporter.service
[Unit]
Description=NetBird API Exporter
Documentation=https://github.com/matanbaruch/netbird-api-exporter
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=netbird-api-exporter
Group=netbird-api-exporter
ExecStart=/usr/local/bin/netbird-api-exporter
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netbird-api-exporter

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes

# Environment
Environment=NETBIRD_API_URL=https://api.netbird.io
Environment=LISTEN_ADDRESS=:8080
Environment=METRICS_PATH=/metrics
Environment=LOG_LEVEL=info

# Set your NetBird API token here or use a separate environment file
# Environment=NETBIRD_API_TOKEN=your_token_here

# Alternative: Load environment from file
# EnvironmentFile=/etc/netbird-api-exporter/config

[Install]
WantedBy=multi-user.target
